# roles/zookeeper/templates/zookeeperd.conf.upstart
#Set the start and stop of the runlevel
start on runlevel [2345]
stop on runlevel [!2345]

#Automatic restart
respawn
respawn limit 2 5

umask 007

kill timeout 30

#Set the start of the user
setuid zookeeper
setgid zookeeper

#Note the use of the start-foreground parameter, it will zookeeper process will run in the foreground, so upstart can control this process
script
exec  {{ pkg_ic.install_path }}/{{ pkg_ic.basename }}/bin/zkServer.sh start-foreground
end script